# byob.service
[Unit]
Description=BYOB Server - Open Source
After=network.target

[Service]
User=REPLACE_USERNAME_HERE
WorkingDirectory=~/byob/web-gui/docker-pyinstaller1
ExecStart=sh ~/byob/web-gui/startup.sh > ~/bootspool.log
ExecStop=/bin/kill -s QUIT $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
