# byob.service
[Unit]
Description=BYOB Server - Open Source
After=network.target

[Service]
User=REPLACE_THIS_USERNAME
WorkingDirectory=/home/REPLACE_THIS_USERNAME/byob/web-gui/docker-pyinstaller1
ExecStart=sh /home/REPLACE_THIS_USERNAME/byob/web-gui/startup.sh > /home/REPLACE_THIS_USERNAME/bootspool.log
ExecStop=/bin/kill -s QUIT $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
