# byob.service
[Unit]
Description=BYOB Server - Open Source
After=network.target

[Service]
User=REPLACE_THIS_USERNAME
WorkingDirectory=$HOME/byob/web-gui/
ExecStart=/usr/bin/sh $HOME/byob/web-gui/service.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
